# qr3_16.conf
[program:qr3_16]
command=python3 /home/curator3/qr3_16/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_16/test.err.log
stdout_logfile=/home/curator3/qr3_16/test.err.log
